# config/includes.chroot/etc/systemd/system/wpa_supplicant@.service.d/exec-override.conf
[Service]
ExecStart=/sbin/wpa_supplicant -c/etc/wpa_supplicant.conf -Dnl80211,wext -i%I